# systemd/subreddit-monitor-updater.timer
[Unit]
Description=Periodically run the subreddit monitor

[Timer]
OnCalendar=*:0/10:0

[Install]
WantedBy=multi-user.target
